# deploy/randomania-backend.service
[Unit]
Description=Randomania backend (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/randomania/backend
Environment="COOKIE_SECURE=true"
ExecStart=/opt/randomania/backend/venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
